# debian/postinst
#!/bin/sh

set -e
sudo setcap 'cap_sys_nice=eip' /usr/bin/modem
sudo systemctl enable 5gmag-rt-modem
sudo systemctl restart 5gmag-rt-modem
